# P2/init.d/eng10032lab12_1
#! /bin/sh
### BEGIN INIT INFO
# Provides:          eng10032lab012_1
# Required-Start:    
# Should-Start:      
# Required-Stop:     
# Default-Start:     S
# Default-Stop:
# Short-Description: Loads led driver
# Description:       Loads led driver
### END INIT INFO

case "$1" in
    start)
    	modprobe eng10032-buzzer
    	;;
    restart|force-reload)
    	rmmod eng10032-buzzer
	modprobe eng10032-buzzer
	;;
    stop)
    	rmmod eng10032-buzzer
	;;
    status)
    	lsmod | grep eng10032-buzzer
	;;
    *)
	echo -n "Usage: $0 " 
	echo "{start|stop|restart|force-reload|status}"
	exit 1
	;;
esac

exit 0

#//DUPLA: Laurien Santin e Luis Miguel Batista
